# g13/systemd/g13-ops.service
[Unit]
Description=G13 Keyboard Driver
Documentation=https://github.com/AreteDriver/G13_Linux
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/g13-ops-gui
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority

# Run as the logged-in user (change if needed)
User=arete
Group=arete

# Permissions for USB HID access
SupplementaryGroups=plugdev

[Install]
WantedBy=graphical-session.target
