# g13/systemd/g13-ops-user.service
# User-level systemd service for G13 driver
# Install to: ~/.config/systemd/user/g13-ops.service
# Enable with: systemctl --user enable g13-ops.service
# Start with: systemctl --user start g13-ops.service

[Unit]
Description=G13 Keyboard Driver (User)
Documentation=https://github.com/AreteDriver/G13_Linux
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/g13-ops-gui
Restart=on-failure
RestartSec=5

# Environment for GUI apps
Environment=QT_QPA_PLATFORM=xcb

[Install]
WantedBy=default.target
